# job-tracker-api.service
[Unit]
Description=Job Tracker API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/job-tracker
ExecStart=/home/ubuntu/job-tracker/venv/bin/python /home/ubuntu/job-tracker/run.py api
Restart=on-failure
RestartSec=5s
Environment="PYTHONPATH=/home/ubuntu/job-tracker"
Environment="PORT=8001"

[Install]
WantedBy=multi-user.target
