# job-tracker-api.service
[Unit]
Description=Job Tracker API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/job-tracker
Environment="PYTHONPATH=/home/ubuntu/job-tracker"
Environment="PORT=8001"
ExecStartPre=/bin/bash -c '/home/ubuntu/job-tracker/venv/bin/python /home/ubuntu/job-tracker/free_port.py 8001'
ExecStart=/home/ubuntu/job-tracker/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8001
ExecStop=/bin/pkill -f "uvicorn main:app"
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=120
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
